# systemd/paipai-memory-ingest.service
[Unit]
Description=派派记忆摄入（扫 Claude session + inbox 到原始湖）
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/paipai
EnvironmentFile=/root/paipai/.env
ExecStart=/usr/bin/python3 /root/paipai/memory_ingest.py
StandardOutput=append:/root/paipai/memory.log
StandardError=append:/root/paipai/memory.log
